# mtconnect/init/mtc_adapter001.conf
description "MTConnect Adapter Simulator"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [016]

exec ruby /etc/mtconnect/adapter/run_scenario.rb -l -p 7878 /etc/mtconnect/adapter/data/mtc_adapter001.txt

respawn
respawn limit 10 600
